# supervisord.conf
[supervisord]

nodaemon=true

[program:myweb]

command=/usr/bin/python main.p 80
